# etc/scollector_supervisord.conf
[program:scollector]
command=#BASE#/scollector-linux -d -h #HEAD#:8070
autostart=true
autorestart=true
stderr_logfile=/var/log/scollector.err.log
stdout_logfile=/var/log/scollector.out.log
